# launchd/com.dailyb.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.dailyb.weekly</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/luke/Service App/DailyB/run_weekly.sh</string>
  </array>

  <!-- 일요일 21:00 KST (시스템 시간대가 Asia/Seoul 기준) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>21</integer>
    <key>Minute</key><integer>0</integer>
    <key>Weekday</key><integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/luke/Service App/DailyB/logs/launchd_weekly.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/luke/Service App/DailyB/logs/launchd_weekly_err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    <key>HOME</key>
    <string>/Users/luke</string>
  </dict>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
